# packaging/systemd/pulse-agent.service
[Unit]
Description=Pulse Agent standalone host service
Documentation=file:/usr/share/doc/pulse-agent/operations.md
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=simple
User=pulse-agent
Group=pulse-agent-admin
SupplementaryGroups=docker
WorkingDirectory=/var/lib/pulse-agent
EnvironmentFile=/etc/pulse-agent/pulse-agent.env
ExecStartPre=/usr/local/bin/pulse-agent config validate --config /etc/pulse-agent/config.json
ExecStart=/usr/local/bin/pulse-agent standalone --config /etc/pulse-agent/config.json
Restart=on-failure
RestartSec=5s
UMask=0007
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/pulse-agent
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
